# arcsight_exporter.service
[Unit]
Description=ArcSight SmartConnector Prometheus Exporter
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/arcsight_exporter/arcsight_exporter.py \
    --base-dir /opt/arcsight/connectors \
    --port 9116 \
    --interval 30 \
    --log-level INFO
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcsight-exporter

[Install]
WantedBy=multi-user.target
